# packer/zjusct-full/rootfs/etc/systemd/system/sssd.service.d/override.conf
[Unit]
# Fix SSSD and docker.socket circular dependency
# Without this, docker.socket depends on sssd.service which depends on docker.socket
DefaultDependencies=no
After=network.target
Before=docker.service autofs.service
