# packer/debian/rootfs/etc/systemd/system/sssd.service.d/override.conf
# 解决 sssd 和 docker.socket 循环依赖
[Unit]
# 1. 核心修复：解除默认依赖限制，允许 SSSD 在早期启动，打破与 basic.target 的循环
DefaultDependencies=no
# 2. 明确启动顺序：确保在 Docker 和 Autofs 启动前，SSSD 已经准备好身份解析服务
# 这样 Docker 启动时就能正确识别来自 IPA 的 docker 组
Before=docker.socket sockets.target autofs.service basic.target
# 3. 补充必要基础：虽然解除了默认依赖，但必须保证底层日志和基础文件系统已就绪
After=systemd-journald.socket local-fs.target
